# umesos-slave/etc/umesos-slave/umesos-cluster-slave.conf
description "umesos Cluster Slave (Mesos Slave)"
author "daniel.kang@ucloud.cn"
start on stopped rc RUNLEVEL=[2345]
respawn
script
  source /etc/umesos-slave/slave_env.sh && mesos-slave
end script
