# umesos-deploy/etc/umesos/service/umesos-cluster-master.conf
description "umesos Cluster Master (Mesos master)"
author "daniel.kang@ucloud.cn"
start on started docker
stop on runlevel [!2345]
respawn
script
  docker-compose -p umesos -f /etc/umesos/compose/umesos-cluster-master.yml up
end script
pre-stop script
  docker-compose -p umesos -f /etc/umesos/compose/umesos-cluster-master.yml stop
end script

